# com.xuxu.a-share-ai-system.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.xuxu.a-share-ai-system.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/xuxu/Documents/Codex/2026-06-23/chrome-plugin-chrome-openai-bundled-file/outputs/a-share-ai-system/scripts/run_local_analysis.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>13</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/xuxu/Library/Logs/a-share-ai-system/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/xuxu/Library/Logs/a-share-ai-system/launchd.err.log</string>
</dict>
</plist>
